# deploy/systemd/polymarket-arbitrage.service
[Unit]
Description=Polymarket Arbitrage Trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/polymarket-arbitrage
EnvironmentFile=/opt/polymarket-arbitrage/.env
ExecStart=/opt/polymarket-arbitrage/venv/bin/python run_auto_trading.py --mode live --continuous --styles up_down --timeframes 5m,15m,1h,4h,1d --assets BTC,ETH,SOL --limit-events 30 --min-position-usdc 1 --max-position-usdc 3 --min-marketable-buy-usdc 1 --max-candidates 1 --scan-interval 10 --enable-auto-claim --claim-interval 300
Restart=always
RestartSec=15
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
